# systemd/dns.service
[Unit]
Description=DNS updater for Namecheap

[Service]
Type=oneshot
ExecStart=/etc/systemd/system/dns.sh

[Install]
WantedBy=suspend.target
