# systemd/dns.timer
[Unit]
Description=update DNS periodically

[Timer]
OnCalendar=*:00/30
Persistent=true

[Install]
WantedBy=timers.target