# systemd/periodic.service
[Unit]
Description=Run a script when periodic.timer triggers

[Service]
ExecStart=/path/to/myscript