# systemd/signald.service
[Unit]
Description=signald

[Service]
ExecStart=/usr/bin/signald
Restart=always

[Install]
WantedBy=multi-user.target
